# supervisor/microros_agent.conf
[program:microros-agent]
command=/home/uwtec/USV/supervisor/microros_agent.sh
user=uwtec
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/uwtec/USV/logs/microros_agent.log
