# supervisor/uwtec_reconnect.conf
[program:uwtec-reconnect]
command=/home/uwtec/USV/supervisor/uwtec_reconnect.sh
user=uwtec
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/uwtec/USV/logs/uwtec_agent.log

