# supervisor/uwtec_agent.conf
[program:uwtec-agent]
command=/home/uwtec/USV/supervisor/uwtec_agent.sh
user=uwtec
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/uwtec/USV/logs/uwtec_agent.log

